# br2-scratchOS/package/applications/template-py/src/main/resources/template-py.service
[Unit]
Description=Process Control System
After=syslog.target

[Service]
ExecStart=/usr/bin/template-py -n
Restart=always

[Install]
WantedBy=multi-user.target
